# collector/launchd/com.alejo.search-bar.collector.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.alejo.search-bar.collector</string>
  <key>ProgramArguments</key>
  <array>
    <string>__RUNNER__</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__DIRECTORY__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>SEARCH_BAR_COLLECTOR_HOST</key>
    <string>__HOST__</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>__LOG_DIRECTORY__/collector.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIRECTORY__/collector-error.log</string>
</dict>
</plist>
